# SSTIBon/var/run/user/1000/systemd/generator.late/xdg-desktop-autostart.target.wants/app-print\x2dapplet@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/etc/xdg/autostart/print-applet.desktop
PartOf=graphical-session.target

Description=Print Queue Applet
After=graphical-session.target

[Service]
Type=exec
ExitType=cgroup
ExecStart=:/usr/bin/system-config-printer-applet
Restart=no
TimeoutSec=5s
Slice=app.slice
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "" "KDE:GNOME:Cinnamon:LXDE:Unity"
