# SSTIBon/var/run/user/1000/systemd/generator.late/xdg-desktop-autostart.target.wants/app-gnome\x2dshell\x2doverrides\x2dmigration@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/etc/xdg/autostart/gnome-shell-overrides-migration.desktop
PartOf=graphical-session.target

Description=GNOME settings overrides migration
After=graphical-session.target

[Service]
Type=exec
ExitType=cgroup
ExecStart=:/usr/libexec/gnome-shell-overrides-migration.sh
Restart=no
TimeoutSec=5s
Slice=app.slice
