# SSTIBon/var/run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-a34a7925\x2df008\x2d4d35\x2d9fa8\x2d089389334bbf.target

[Mount]
What=/dev/disk/by-uuid/a34a7925-f008-4d35-9fa8-089389334bbf
Where=/
Type=ext4
Options=errors=remount-ro
