# SSTIBon/var/run/systemd/generator/boot-efi.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2duuid-012A\x2dFCBD.service
After=systemd-fsck@dev-disk-by\x2duuid-012A\x2dFCBD.service
After=blockdev@dev-disk-by\x2duuid-012A\x2dFCBD.target

[Mount]
What=/dev/disk/by-uuid/012A-FCBD
Where=/boot/efi
Type=vfat
Options=umask=0077
