# SSTIBon/var/run/user/1000/systemd/generator.late/app-ubuntu\x2dreport\x2don\x2dupgrade@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/etc/xdg/autostart/ubuntu-report-on-upgrade.desktop
PartOf=graphical-session.target

Description=Ubuntu report try to sends metrics data on release upgrade
After=graphical-session.target

[Service]
Type=exec
ExitType=cgroup
ExecStart=:/usr/bin/ubuntu-report send upgrade
Restart=no
TimeoutSec=5s
Slice=app.slice
